# joybox.service
[Unit]
Description=Joybox

[Service]
ExecStart=/home/pi/joybox/start.sh

[Install]
WantedBy=default.target